# scripts/omero-quota-enforcer.path
[Unit]
Description=Watch for OMERO quota updates

[Path]
PathModified=/etc/default/omero-quota-enforcer
# The actual path will be injected by the installer

[Install]
WantedBy=multi-user.target